# nm-applet.service
[Unit]
Description=nm-applet

[Service]
ExecStart=/usr/bin/nm-applet

[Install]
WantedBy=default.target
